# deploy/victoria-metrics.service
[Unit]
Description=VictoriaMetrics — time-series metrics storage
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin

# binary lives at /home/admin/victoria-metrics-prod (see deploy-vm.sh)
ExecStart=/home/admin/victoria-metrics-prod \
    -storageDataPath=/home/admin/vm-data \
    -retentionPeriod=3   \
    -httpListenAddr=:8428

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=victoria-metrics

[Install]
WantedBy=multi-user.target
